# web/recorder-web.service
[Unit]
Description=Recorder admin web (Flask + gunicorn)
After=network.target

[Service]
Type=simple
User=ftadmin
WorkingDirectory=/opt/recorder/web
Environment="RECORDINGS_DIR=/opt/recorder/recordings"
Environment="SRS_PORT=8080"
Environment="STREAM_MAIN=recorder-main"
Environment="STREAM_AUX=recorder-aux"
Environment="SESSION_TIMEOUT_MIN=30"
# --timeout 300: LLM endpoints (transcript/refine, summary) call out to
#   external LLM APIs that routinely take 60-180s on long meetings.
#   Default gunicorn worker timeout is 30s — it kills the worker mid-call
#   and the browser sees "Failed to fetch".
# --graceful-timeout 30: on systemctl reload, still let workers drain in 30s.
# --bind 127.0.0.1:8088: 仅本机监听。生产必须前置 nginx/caddy 反代终结 TLS。
#   旧值 0.0.0.0:8088 让 session cookie 走明文 HTTP，等价于绕过登录。
#   反代模板 + 防火墙配置见 docs/deployment.md "HTTPS 反代" 章节。
# --forwarded-allow-ips 127.0.0.1: 信任本机反代发来的 X-Forwarded-Proto / -For 头。
ExecStart=/usr/bin/python3 -m gunicorn --workers 2 --threads 4 --timeout 300 --graceful-timeout 30 --bind 127.0.0.1:8088 --forwarded-allow-ips 127.0.0.1 --access-logfile - --error-logfile - app:app
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
